# deploy/adsb-decode.service
[Unit]
Description=adsb-decode dashboard
After=network.target

[Service]
Type=simple
User=adsb
WorkingDirectory=/opt/adsb-decode
ExecStart=/opt/adsb-decode/adsb serve --db-path /opt/adsb-decode/data/adsb.db --port 8000 --host 127.0.0.1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
