# deploy/adsb-decode-healthcheck.service
[Unit]
Description=adsb-decode API healthcheck + auto-restart
After=adsb-decode.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/adsb-decode-healthcheck.sh
