# deploy/adsb-decode.service
[Unit]
Description=adsb-decode tracking dashboard
After=network.target postgresql.service

[Service]
Type=simple
User=adsb
WorkingDirectory=/opt/adsb-decode
EnvironmentFile=/opt/adsb-decode/.env
ExecStart=/opt/adsb-decode/adsb serve --port 8000 --host 127.0.0.1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
